# contrib/OpenRC/init_d
#!/sbin/openrc-run
#
# /etc/init.d/automua

: ${AUTOMUA_CONF:="/etc/${RC_SVCNAME}.conf"}
: ${AUTOMUA_USER:="automua"}
: ${AUTOMUA_ARGS:="--port 4243"}

command="/usr/bin/python"
command_args="/usr/bin/flask run ${AUTOMUA_ARGS}"
command_background="true"
command_user="${AUTOMUA_USER}"
pidfile="/run/${RC_SVCNAME}.pid"
required_files="${AUTOMUA_CONF}"

depend() {
	use logger net
	before nginx
}

start_pre() {
	export AUTOMUA_CONF
	export EPYTHON="python3.9"
	export FLASK_APP="automua.server:app"
	export FLASK_ENV="production"
}
